# com.sticks3.voice-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  StickS3 voice-recorder daily summary+notify at 22:00 (local time).

  Install:
    1. Edit the two ABSOLUTE paths below to match this checkout.
    2. cp com.sticks3.voice-daily.plist ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/com.sticks3.voice-daily.plist
  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.sticks3.voice-daily.plist
  Test now:
    launchctl start com.sticks3.voice-daily
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sticks3.voice-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/path/to/Taskhub-for-StickS3/run_daily.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/path/to/Taskhub-for-StickS3/daily.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/Taskhub-for-StickS3/daily.log</string>
</dict>
</plist>
