# com.sticks3.transcribe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Keep the transcription worker (transcribe.py) running, restart on crash/boot.
  Depends on the whisper-server LaunchAgent being up (WHISPER_SERVER_URL).

  Install:
    cp com.sticks3.transcribe.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.sticks3.transcribe.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sticks3.transcribe</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/path/to/Taskhub-for-StickS3/run_transcribe.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/path/to/Taskhub-for-StickS3/transcribe.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/Taskhub-for-StickS3/transcribe.log</string>
</dict>
</plist>
